# lxc/container.timer
#
# ~/.local/share/systemd/user/container.timer
#

[Unit]
Description=Timer to run container

[Timer]
OnBootSec=30

[Install]
WantedBy=timers.target